# daemon/aztec-butler.service
[Unit]
Description=Aztec Butler - Prometheus monitoring and automation for Aztec nodes
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aztec-butler
ExecStart=/usr/bin/node dist/index.js serve
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin:/bin

[Install]
WantedBy=multi-user.target
